# systemd/ohff-teletext.timer
[Unit]
Description=Regenerate OHFF teletext EP1 page every 5 minutes

[Timer]
OnBootSec=1min
OnUnitActiveSec=5min
AccuracySec=30s
Unit=ohff-teletext.service

[Install]
WantedBy=timers.target
